# root/etc/services.d/cron/run
#!/usr/bin/execlineb -P

with-contenv
crond -s -n -P $PATH
